# cron/cron.d/finance-scraper
5 9 * * 1-5 root cd /home/nagiyu/finance/containers && docker compose -f docker-compose-scraper.yml build --no-cache && docker compose -f docker-compose-scraper.yml -p finance up -d
0 15 * * 1-5 root cd /home/nagiyu/finance/containers && docker compose -f docker-compose-scraper.yml -p finance down
55 16 * * 1-5 root cd /home/nagiyu/finance/containers && docker compose -f docker-compose-scraper.yml build --no-cache && docker compose -f docker-compose-scraper.yml -p finance up -d
0 9 * * 2-6 root cd /home/nagiyu/finance/containers && docker compose -f docker-compose-scraper.yml -p finance down
