# cron/cron.d/letsencrypt-renew
0 0 1 * * root certbot renew --pre-hook "systemctl stop nginx" --post-hook "systemctl start nginx"
